# bitrix/crontab/crontab.cfg
0 */1 * * * php -f /var/www/html/bitrix/php_interface/include/catalog_import/cron_frame.php 1 >/var/www/html/bitrix/php_interface/include/catalog_import/logs/1.txt
